# roles/arr/files/etc/systemd/system/docker.nzbget.service
[Unit]
Description=nzbget on Docker
Requires=docker.service network.target mnt-nzbget.mount

[Service]
User=nzb
Group=nzb
Restart=always
ExecStart=/usr/bin/docker-compose -f /srv/docker/nzbget/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /srv/docker/nzbget/docker-compose.yml down

[Install]
WantedBy=default.target
